# claude_vid.service
[Unit]
Description=LightBox Next.js Frontend (claude-vid)
Documentation=https://vid.arc-codex.com
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/mnt/data/claude_vid

ExecStart=/mnt/data/claude_vid/claude_vid.sh start
ExecStop=/mnt/data/claude_vid/claude_vid.sh stop

# Allow up to 3 minutes for Docker build/start on first boot
TimeoutStartSec=180
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-vid

[Install]
WantedBy=multi-user.target
